# docker/supervisor/etc/conf.d/httpd.ini
[program:httpd]
command=/usr/local/bin/%(program_name)s.sh
stopsignal=WINCH
stopwaitsecs=10
