# docker/supervisor/etc/conf.d/pythonserver.ini
[program:pythonserver]
command=/usr/local/bin/%(program_name)s.sh
